# install/autostart/grunt.service
# grunt.service — auto-boot the grunt robot stack
#
# Install:
#   sudo cp /home/karim/ros2_ws/src/grunt/install/autostart/grunt.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable grunt.service     # auto-start on every boot
#   sudo systemctl start grunt.service      # start right now
#
# Watch live output (color + full detail, like attaching to a terminal):
#   journalctl -fu grunt
#
# Stop / restart:
#   sudo systemctl stop grunt
#   sudo systemctl restart grunt
#
# Switch active site (cascades to missions_dir for all consumers):
#   sudo grunt-set-site dallas_makerspace
#   grunt-set-site --show
# (Or manually: edit /etc/default/grunt, restart grunt.service.)
#
# Shutdown path: ExecStop parks the arm to the bumper preset and
# publishes zero cmd_vel briefly. systemd then SIGTERMs the whole
# cgroup (KillMode=control-group) — including p2os_driver, whose
# rclcpp signal handler runs Shutdown() (SONAR-off → STOP → CLOSE)
# before the process exits. TimeoutStopSec is sized for Nav2
# lifecycle deactivation + the driver's serial writes to complete.

[Unit]
Description=Grunt robot autonomy stack (ROS 2 Jazzy)
After=network-online.target
Wants=network-online.target

# Pull in Nav2 (indoor by default) when grunt.service starts manually.
# `systemctl enable grunt-nav@indoor.service` is only honored at boot
# via multi-user.target — a mid-session `systemctl start grunt.service`
# won't activate it without this Wants. BindsTo=grunt.service on the
# nav unit already handles the reverse (nav stops when grunt stops).
#
# Caveat: each `systemctl restart grunt.service` bounces Nav2 back to
# indoor even if the operator had switched to outdoor earlier. Any
# runtime mode change is discarded on full-stack restart; the
# nav_mode_switcher can be used to flip modes again post-start.
#
# To change the default mode, edit this line to point to
# grunt-nav@outdoor.service.
#
# Ordering: deliberately NO After= here. The nav unit already has
# After=grunt.service, so nav waits for grunt. Adding the reverse
# would be circular.
Wants=grunt-nav@indoor.service

[Service]
Type=simple
User=karim
Group=karim
Environment=GRUNT_PREFIX=grunt1
Environment=GRUNT_SITE=ranchero
# NTRIP caster mountpoint for RTCM corrections. Override per site
# in /etc/default/grunt — eventually replaced by a site→mountpoint
# map (PRD: ponderdocs/nav/docs/prd_ntrip_mountpoint_switching.md).
Environment=GRUNT_NTRIP_MOUNTPOINT=KubotaCentral
Environment=HOME=/home/karim
WorkingDirectory=/home/karim/ros2_ws

# Operator-controlled overrides live in /etc/default/grunt (standard
# Linux pattern). Currently used to switch the active site:
#   echo 'GRUNT_SITE=dallas_makerspace' | sudo tee /etc/default/grunt
#   sudo systemctl restart grunt
# Or use the grunt-set-site helper script for the same effect.
# Leading dash on the path: missing file is fine, fall back to the
# Environment= defaults above.
EnvironmentFile=-/etc/default/grunt

# exec replaces the bash shell so systemd tracks the ros2 launch PID
# directly. sourcing is inline because systemd's Environment= doesn't
# run shell expansions.
ExecStart=/bin/bash -c 'source /opt/ros/jazzy/setup.bash && source /home/karim/ros2_ws/install/setup.bash && exec ros2 launch grunt_bringup grunt.quiet.launch.py prefix:=${GRUNT_PREFIX} site:=${GRUNT_SITE} mountpoint:=${GRUNT_NTRIP_MOUNTPOINT} use_sonar:=True'

ExecStop=/home/karim/ros2_ws/src/grunt/install/autostart/grunt_park_and_stop.sh ${GRUNT_PREFIX}

# 25s covers: ExecStop (~5s) + Nav2 lifecycle shutdown (5-10s) +
# p2os Shutdown() serial writes (~1s). Short enough that a wedged
# stop doesn't delay reboot indefinitely.
TimeoutStopSec=25

# Disable per-service CPU accounting so the kernel CFS scheduler
# distributes CPU across individual processes instead of across
# cgroups. Observed: splitting Nav2 into its own systemd unit
# (grunt-nav@.service) halved Nav2's effective CPU share on the
# quad-core NUC, causing MPPI trajectory stutter. With accounting
# off, both cgroups compete process-by-process as before.
# Same setting lives on grunt-nav@.service.
CPUAccounting=no

# Restart on crash; `systemctl stop` stays stopped.
Restart=on-failure
RestartSec=10

# Cgroup-wide SIGTERM is the critical bit — every node's own rclcpp
# signal handler runs, not just the launch tool's.
KillMode=control-group
KillSignal=SIGTERM

[Install]
WantedBy=multi-user.target
