# grunt_bringup/systemd/grunt-wifi-on.service
[Unit]
Description=Force WiFi radio on at boot (headless-robot safety net)
After=NetworkManager.service
Wants=NetworkManager.service

[Service]
Type=oneshot
# Wait for NM startup (-s = ready for D-Bus commands, NOT for a connection).
# Without -s this blocks the full timeout when wifi is off -- the exact case we
# need to recover from. Leading '-' makes a non-zero exit non-fatal so the
# radio-on below always runs regardless.
ExecStart=-/usr/bin/nm-online -q -s -t 30
ExecStart=/usr/bin/nmcli radio wifi on
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
